# build/ota-update-server/debian/ota.service
[Unit]
Description=OTA Update Server
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/ota /root/key.txt
PrivateTmp=true

[Install]
WantedBy=multi-user.target
