# fett/cyberPhys/infotainment-server/debian/infotainment-server.service
[Unit]
Description=Infotainment Server
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/infotainment_server 10.88.88.255
PrivateTmp=true

[Install]
WantedBy=multi-user.target
